# systemd/cross-server-worker.service
[Unit]
Description=Cross-Server Copilot Worker — Background task executor
After=network.target xray-vpn.service
Wants=xray-vpn.service

[Service]
Type=simple
ExecStart=/bin/bash /root/.copilot/cross-server/worker.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cross-server-worker
Environment=HOME=/root
WorkingDirectory=/root

# Security: restrict dangerous syscalls
NoNewPrivileges=no
ProtectSystem=no

[Install]
WantedBy=multi-user.target
